# crontab
0 0 7 1/1 * ? * /usr/bin/python /home/pi/AutoHBL/main.py > /home/AutoHBL/main.log 2>&1